# packaging/slackbuild/fono/fono.service
[Unit]
Description=Fono voice dictation daemon
Documentation=https://github.com/bogdanr/fono
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/fono
Restart=on-failure
RestartSec=3

# Security hardening (this is a per-user unit, so minimal caps needed).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/fono %h/.local/share/fono %h/.cache/fono %h/.local/state/fono
PrivateTmp=true

[Install]
WantedBy=default.target
